# rpisys/scripts/crontab
SHELL=/bin/bash
BASH_ENV=$HOME/.bashrc
@reboot sleep 5 && tmux new -d -s ros-session $HOME/ros/startup.sh

